# scripts/transformer-monitor.service
[Unit]
Description=Transformer Thermal Monitor
After=network.target network-online.target time-sync.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/transformer-monitor
# Uses the start.sh script which handles environment setup
ExecStart=/home/pi/transformer-monitor/scripts/start.sh
Restart=always
RestartSec=10
# Standard output logging to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=transformer-monitor

[Install]
WantedBy=multi-user.target
